# etc/service/heavy-logging/run
#!/bin/bash -e

while true; do
  date
  echo -e 'first\nsecond\nthird\nline'
done
