# etc/service/heavy-logging/run
#!/bin/bash -e

while true; do
  date
done
